# Монтана_Montana_蒙大拿/Русский/бот/junona.service
# /etc/systemd/system/junona.service
# Montana Protocol — Юнона Bot Service
#
# Установка на каждом узле:
# 1. Скопировать файл: sudo cp junona.service /etc/systemd/system/
# 2. Изменить MONTANA_NODE_NAME на имя узла (amsterdam/moscow/almaty/spb/novosibirsk)
# 3. Reload: sudo systemctl daemon-reload
# 4. Enable: sudo systemctl enable junona
# 5. Start: sudo systemctl start junona
# 6. Logs: journalctl -u junona -f

[Unit]
Description=Junona Montana Protocol Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bot
ExecStart=/usr/bin/python3 junomontanaagibot.py

# ═══════════════════════════════════════════════════════════════════════════════
# ВАЖНО: Установить имя узла!
# ═══════════════════════════════════════════════════════════════════════════════
# amsterdam = priority 1 (master)
# moscow    = priority 2
# almaty    = priority 3
# spb       = priority 4
# novosibirsk = priority 5
Environment="MONTANA_NODE_NAME=amsterdam"

# Telegram token из keychain или .env
Environment="TELEGRAM_TOKEN_JUNONA="
Environment="OPENAI_API_KEY="

# Автоперезапуск при падении
Restart=always
RestartSec=5

# Логирование
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
